# docs/deploy/beer-lab-sync.service
# Beer-Lab-Ware sync daemon — systemd unit template
#
# Runs the PRE-BUILT single-file bundle (dist/sync-server.mjs — see "Building
# the daemon" in README.md) directly under plain node. No tsx, no `npm install`
# at runtime, no repo `node_modules` on the serving path — just Node 22+ and the
# bundle. Binds 127.0.0.1 only; your reverse proxy is the sole ingress. Secrets
# (device-token hashes, brewery file path) come from the 0600 EnvironmentFile —
# never the unit or the repo.
#
# Build once per deploy from a checkout (see README.md "Bring-up order"):
#   cd /srv/beer-lab-ware/app && npm ci && npm run build:daemon
# then copy/symlink WorkingDirectory's dist/ (or point WorkingDirectory straight
# at the checkout, as below) so ExecStart's relative path resolves.

[Unit]
Description=Beer-Lab-Ware sync daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=beerlab
Group=beerlab
WorkingDirectory=/srv/beer-lab-ware/app
# EnvironmentFile supplies BREWERY_FILE, SYNC_TOKEN_HASHES, SYNC_PORT,
# SYNC_KEEP_GENERATIONS (see sync.env.example).
EnvironmentFile=/etc/beer-lab-ware/sync.env
ExecStart=/usr/bin/node dist/sync-server.mjs
Restart=on-failure
RestartSec=3
# Hardening: the daemon needs only its brewery.json dir writable.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/beer-lab-ware
# The canonical brewery.json lives here — keep it on storage that survives
# reboots/redeploys of this host.

[Install]
WantedBy=multi-user.target
